# notification_service/config/crontab
TZ=Asia/Taipei
* * * * * root /usr/local/bin/python3 /app/notification_service/main.py >> /var/log/cron.log 2>&1

