# crawler_service/config/crontab
TZ=Asia/Taipei
*/5 * * * * root /usr/bin/timeout 300 /usr/local/bin/python3 /app/crawler_service/main.py >> /var/log/cron.log 2>&1

