# crawler_service/config/crontab
TZ=Asia/Taipei
*/5 * * * * root cd /app/crawler_service && /usr/local/bin/scrapy crawl search_product >> /var/log/cron.log 2>&1
